# config/local.d/rfkill.start
#
# disable bluetooth adapter on startup
#
#rmmod radeon
#rfkill block 1
